# _disabled/hdd-automount/etc/hotplug.d/block/99-hdd
#!/bin/sh
/usr/sbin/manage_hdd.sh "$ACTION" "$DEVNAME"
