# files/etc/init.d/tailscale
#!/bin/sh /etc/rc.common
#
# Minimal OpenWrt procd init unit for tailscaled. Mirrors what the OpenWrt
# 24.10 tailscale package will install once the next IB bake includes the
# package — this is the bg852 HOT-FIX equivalent (manual install on a unit
# whose firmware predates the tailscale-in-PACKAGES change).
#
# After the official package lands in the bake, this file should be REMOVED
# from the device's overlay so the package-managed unit takes over (path
# identical, semantics identical, but maintained by the OpenWrt package).
#

USE_PROCD=1
START=80
STOP=10

PROG=/usr/sbin/tailscaled
STATE_DIR=/var/lib/tailscale
SOCKET=/var/run/tailscale/tailscaled.sock
LOG_FILE=/var/log/tailscaled.log

start_service() {
    [ -x "$PROG" ] || { echo "tailscale: $PROG missing" >&2; return 1; }
    mkdir -p "$STATE_DIR" /var/run/tailscale
    procd_open_instance
    procd_set_param command "$PROG" \
        --state="$STATE_DIR/tailscaled.state" \
        --socket="$SOCKET" \
        --port=41641 \
        --tun=userspace-networking
    procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-0}
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param file "$PROG"
    procd_close_instance
}

stop_service() {
    # nvrd-core's ensureTailscaleUp must NOT bring the daemon down — the marker
    # gates re-bootstrap, but state-dir wipes still re-register on next boot.
    # Stopping the service here is the normal init.d/tailscale stop path
    # (operator-triggered or system shutdown); it does NOT delete state.
    :
}

service_triggers() {
    procd_add_reload_trigger "tailscaled"
}
